# infra/sla-watchdog/digilist-watchdog.timer
[Unit]
Description=Run the Digilist SLA watchdog every 60s

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5s
Unit=digilist-watchdog.service

[Install]
WantedBy=timers.target
